# CLI-Install/ImageStoreFRONT.service
[Service]
Type=simple
Restart=on-failure
WorkingDirectory=/etc/imagestore/frontend
ExecStart=npm start
