# CLI-Install/ImageStoreBACK.service
[Service]
Type=simple
Restart=on-failure
WorkingDirectory=/etc/imagestore/backend
ExecStart=npm start
